# testing/docker/install-agent.service
[Unit]
Description=install netuitive-agent

[Service]
ExecStart=/vagrant/docker/test.sh

[Install]
WantedBy=multi-user.target