# testing/install-agent.service
[Unit]
Description=install netuitive-agent

[Service]
ExecStart=/vagrant/test.sh

[Install]
WantedBy=multi-user.target