# netuitive/upstart/netuitive-agent.conf
# netuitive-agent - A system statistics collector for Netuitive
#
# Netuitve-agent is a daemon and toolset for gather system statistics
# and publishing them to Netuitive.

start on stopped rc RUNLEVEL=[2345]
stop on runlevel [!2345]

respawn

script
    export SENSORS_LIB=/opt/netuitive-agent/embedded/lib/libsensors.so
    /opt/netuitive-agent/bin/supervisord -n -c /opt/netuitive-agent/conf/supervisor.conf
end script

